# templates/default/systemd/kube-proxy.service.erb
[Unit]
Description=kube-proxy
After=network.target

[Service]
ExecStart=/usr/sbin/kube-proxy
