# templates/default/systemd/kubelet.service.erb
[Unit]
Description=kubelet
After=network.target

[Service]
# User=kubernetes
ExecStart=<%= @kubelet_command %>
