# templates/default/systemd/kube-controller-manager.service.erb
[Unit]
Description=kube-controller-manager
After=network.target

[Service]
Type=simple
User=kubernetes
ExecStart=/usr/sbin/kube-controller-manager
