# templates/systemd/kube-scheduler.service.erb
[Unit]
Description=kube-scheduler
After=network.target

[Service]
User=kubernetes
ExecStart=<%= @kube_scheduler_command %>
