# templates/systemd/kube-proxy.service.erb
[Unit]
Description=kube-proxy
After=network.target

[Service]
ExecStart=<%= @kube_proxy_command %>
