# templates/systemd/kubelet.service.erb
[Unit]
Description=kubelet
Documentation=https://k8s.io
After=network.target <%= @container_runtime_service %>
Wants=<%= @container_runtime_service %>

[Service]
# User=kubernetes
ExecStart=<%= @kubelet_command %>

[Install]
WantedBy=multi-user.target
