# templates/systemd/kube-apiserver.service.erb
[Unit]
Description=kube-apiserver
After=network.target

[Service]
Type=notify
User=kubernetes
ExecStart=<%= @kube_apiserver_command %>
