# templates/systemd/kube-apiserver.service.erb
[Unit]
Description=kube-apiserver
After=network.target

[Service]
Type=notify
User=kubernetes
ExecStart=/usr/sbin/kube-apiserver --service-cluster-ip-range=10.0.0.1/24 --etcd-servers=http://127.0.0.1:4001
